# reposync-watcher.service
[Unit]
Description=RepoSync Watcher (inotify)
After=graphical-session.target

[Service]
Type=simple
Environment=DEBOUNCE_MS=500
WorkingDirectory=%h/Repos/Meus/reposync
ExecStart=%h/Repos/Meus/reposync/watch_repos.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target